# runpod/supervisord.conf
; Supervisor config for the RunPod all-in-one image.
;
; A RunPod Pod runs exactly ONE container, so we run two processes under
; supervisord:
;   * engine  – llama.cpp llama-server on 127.0.0.1:8001 (NOT exposed)
;   * gateway – the FastAPI gateway on 0.0.0.0:8000 (the ONLY exposed port)
;
; Clients reach only the authenticated gateway via the RunPod proxy URL; the
; engine is never directly reachable from the internet.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:engine]
command=/engine/entrypoint.sh
autostart=true
autorestart=true
startretries=3
startsecs=0
stopwaitsecs=30
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gateway]
; Port comes from the GATEWAY_PORT env var (the Dockerfile sets it = 8000, the
; same value EXPOSEd and health-checked). Override all three together to change it.
command=/opt/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port %(ENV_GATEWAY_PORT)s --no-access-log --proxy-headers --forwarded-allow-ips *
directory=/app
autostart=true
autorestart=true
startsecs=3
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
